# roles/app_wordpress_stack/templates/redis.container.j2
[Container]
ContainerName=wp_redis
Image=docker.io/library/redis:7-alpine
Pod=wp_stack.pod
